# container/crontabs/news
# min  hour  day  month  weekday  command
  */1    *    *     *       *     /usr/local/sbin/check-health
  */5    *    *     *       *     /usr/local/sbin/extract
